# bin/weathermonitor.service
[Unit]
Description=Run weather recording service

[Service]
Type=simple
ExecStart=/home/swoboj/weathercheck/bin/weathermonitor.sh
Restart=no
WorkingDirectory=/home/swoboj/

[Install]
WantedBy=default.target
