# bin/weatherrecord.service
[Unit]
Description=Run weather recording service

[Service]
Type=simple
ExecStart=/home/swoboj/weathercheck/bin/weatherservice.sh
Restart=always
WorkingDirectory=/home/swoboj/

[Install]
WantedBy=default.target
